# systemd/user/rpi-git-handoff.service
[Unit]
Description=RPi Git handoff sync to Milhy-PC CI gateway

[Service]
Type=oneshot
WorkingDirectory=/home/milhy777/rpi-dashboard
ExecStart=/home/milhy777/rpi-dashboard/tools/trigger-ci-handoff.sh
